# autobump.service
[Unit]
Description=Discord Auto Bump Bot Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=autobump
WorkingDirectory=/home/autobump/auto-bump
ExecStart=/home/autobump/auto-bump/venv/bin/python3 /home/autobump/auto-bump/bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
